# container/conf/fluentd.conf
[program:fluentd]
command=/usr/bin/fluentd
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/fluentd.error.log
stdout_logfile=/var/log/supervisor/fluentd.output.log
